# development/installation/services/sdp_file_verification/sdp-file-verification.path
[Unit]
Description=Anna - Monitoring "/var/run/sdp" file.

[Install]
RequiredBy=bluetooth.service

[Path]
# Requires that "/var/run/sdp" exists before starting the unit.
PathExists=/var/run/sdp
# Unit to execute qhen "/var/run/sdp" is found.
Unit=sdp-file-group-change-verification.service
