# release/0.1/installation/services/sdp_file_verification/sdp-file-verification.path
[Unit]
Description=Monitor \"/var/run/sdp\"

[Install]
WantedBy=bluetooth.service

[Path]
PathExists=/var/run/sdp
Unit=sdp-file-group-change-verification.service
