# installation/units/cleanup/unit_files/cleanup-anna.service
[Unit]
Description=Clean up Anna directories
After=?

[Service]
Type=simple
ExecStart=clean_up_directories.sh
Restart=?

[Install]
WantedBy=multi-user.target
